# config_files/supervisor_config/gunicorn.conf
[program:gunicorn]

command=gunicorn_django --pid /tmp/gunicorn.pid

directory=/var/www/datamining/datamining

user=www-data

autostart=True

autorestart=True

redirect_stderr=True